# service/set-gpu-pl.service
[Service]
ExecStart=nvidia-smi -pl 256

[Install]
WantedBy=multi-user.target
